# crontab
* * * * * /etc/scripts/log-cron.sh
* * * * * /etc/scripts/request-app.sh
